# packaging/deb/postinst
#!/bin/sh
set -e

# Update man database
if command -v mandb >/dev/null 2>&1; then
    mandb -q 2>/dev/null || true
fi

# Install Python dependencies for cloud backup tools from the packaged
# requirements file so source dependency updates and package installs stay aligned.
REQ_FILE=/usr/share/global-sys-utils/requirements.txt

if command -v pip3 >/dev/null 2>&1; then
    echo "Installing Python dependencies for global-logrotate cloud tools..."
    if [ -f "$REQ_FILE" ]; then
        pip3 install --quiet --break-system-packages -r "$REQ_FILE" 2>/dev/null || \
        pip3 install --quiet -r "$REQ_FILE" 2>/dev/null || \
        echo "Warning: pip install failed. Run manually: pip3 install -r $REQ_FILE"
    else
        echo "Warning: $REQ_FILE not found. Install cloud tool dependencies manually."
    fi
else
    echo "Warning: pip3 not found. Install Python dependencies manually:"
    echo "  pip3 install -r $REQ_FILE"
fi

# Reload systemd and enable the daemon service
if command -v systemctl >/dev/null 2>&1 && systemctl is-system-running --quiet 2>/dev/null; then
    systemctl daemon-reload || true
    systemctl enable global-logrotate.service 2>/dev/null || true
    echo ""
    echo "global-logrotate installed."
    echo ""
    echo "  Option A — long-running daemon (recommended, includes disk monitoring):"
    echo "    1. Add SCHEDULE to /etc/global-sys-utils/global.conf or conf.d/"
    echo "    2. systemctl start global-logrotate"
    echo ""
    echo "  Option B — systemd timer (schedule managed by systemd):"
    echo "    systemctl enable --now global-logrotate-once.timer"
    echo ""
fi

exit 0
